# templates/default/systemd.service.erb
<% prefix = @authbind ? '/usr/bin/authbind --deep ' : '' %>
[Unit]
Description = GlassFish Server for <%= @new_resource.domain_name %>
After = syslog.target network.target remote-fs.target

[Service]
User=<%= node['glassfish']['user'] %>
ExecStart = <%= "#{prefix} #{@start_domain_command}" %>
ExecStop = <%= "#{prefix} #{@stop_domain_command}" %>
ExecReload = <%= "#{prefix} #{@restart_domain_command}" %>
Type = forking
<%- if @new_resource.env_variables
  @new_resource.env_variables.sort.each do |key, value| %>
Environment="<%= key %>=<%= value %>"
<% end -%>
<% end -%>
<%- if @new_resource.env_var_file %>
EnvironmentFile=<%= @new_resource.env_var_file %>
<% end -%>
TimeoutStartSec = <%= @start_domain_timeout %>
TimeoutStopSec = <%= @stop_domain_timeout %>

[Install]
WantedBy = multi-user.target
